# docker/service/tail/run
#!/bin/bash
sleep 2
exec tail -F /docker/service/{rails,sidekiq}/log/current
